# docker/casmgr/s6/s6-rc.d/20-postgres/finish
#!/usr/bin/with-contenv bash
# Executed by s6-supervise once the `run` script's process has exited. The
# TERM trap in `run` already performs a native Postgres fast shutdown; this
# is a defensive no-op safety net (e.g. if the process died without the trap
# firing, such as an OOM-kill) so the next boot never finds a stale postmaster
# lock or needs WAL crash-recovery replay.
gosu postgres pg_ctl stop -D /var/casmgr/postgres/data -m fast >/dev/null 2>&1 || true
exit 0
