# packaging/systemd/casmgr.service
[Unit]
Description=casmgr (all-in-one container: PostgreSQL + code_analysis daemon)
Documentation=man:casmgr-server(8) info:casmgr-server
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/casmgr
ExecStartPre=-/bin/sh -c 'docker network inspect smart-assistant >/dev/null 2>&1 || docker network create smart-assistant'
ExecStartPre=-/usr/bin/docker compose -f /etc/casmgr/docker-compose.yml pull
ExecStart=/usr/bin/docker compose -f /etc/casmgr/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /etc/casmgr/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
