# debian/postinst
#!/bin/bash
# debian/postinst — casuser, directories, image provisioning, systemd.

set -e

CASMGR_USER="casuser"
CASMGR_GROUP="casgrp"
INSTALL_DIR="/usr/lib/casmgr-server"
CONFIG_DIR="/etc/casmgr"
DATA_DIR="/var/casmgr"
LOG_DIR="/var/log/casmgr"
SECRETS="${DATA_DIR}/secrets"
CONFIG="${CONFIG_DIR}/config.json"
CONFIG_TEMPLATE_LOCAL="${CONFIG_DIR}/config.json.template"
RUNTIME_CONFIG_TEMPLATE="/usr/share/casmgr-server/config.json.template"
BIN="/usr/lib/casmgr/bin"
COMPOSE_FILE="${CONFIG_DIR}/docker-compose.yml"

ensure_runtime_config_template() {
    if [[ -f "$RUNTIME_CONFIG_TEMPLATE" ]]; then
        return 0
    fi
    install -d -m 0755 /usr/share/casmgr-server
    if [[ -f /usr/share/doc/casmgr-server/config.json.template ]]; then
        install -m 0644 /usr/share/doc/casmgr-server/config.json.template \
            "$RUNTIME_CONFIG_TEMPLATE"
        return 0
    fi
    if [[ -f /usr/share/doc/casmgr-server/config.json.template.gz ]]; then
        zcat /usr/share/doc/casmgr-server/config.json.template.gz >"$RUNTIME_CONFIG_TEMPLATE"
        chmod 0644 "$RUNTIME_CONFIG_TEMPLATE"
        return 0
    fi
    echo "ERROR: config template missing (expected ${RUNTIME_CONFIG_TEMPLATE} or doc copy)" >&2
    return 1
}

seed_watch_catalog_example() {
    local example_src="/usr/share/casmgr-server/watch-catalog-example"
    local catalog="${DATA_DIR}/watch_catalog"
    [[ -d "$example_src" ]] || return 0
    local uuid_dir
    for uuid_dir in "$example_src"/*/; do
        [[ -d "$uuid_dir" ]] || continue
        local uuid
        uuid="$(basename "$uuid_dir")"
        local dest="${catalog}/${uuid}"
        if [[ -e "$dest" ]]; then
            continue
        fi
        cp -a "$uuid_dir" "$dest"
        chown -R "$CASMGR_USER:$CASMGR_GROUP" "$dest"
        chmod 0750 "$dest"
        if [[ -f "${dest}/settings.json" ]]; then
            chmod 0640 "${dest}/settings.json"
        fi
    done
}

create_user() {
    if ! getent group "$CASMGR_GROUP" >/dev/null; then
        addgroup --system "$CASMGR_GROUP"
    fi
    if ! id "$CASMGR_USER" >/dev/null 2>&1; then
        adduser --system --ingroup "$CASMGR_GROUP" --no-create-home --shell /usr/sbin/nologin "$CASMGR_USER"
    fi
}

create_layout() {
    install -d -m 0750 -o root -g "$CASMGR_GROUP" "$CONFIG_DIR" "${CONFIG_DIR}/mtls"
    install -d -m 0750 -o "$CASMGR_USER" -g "$CASMGR_GROUP" \
        "$DATA_DIR" "$LOG_DIR" \
        "${DATA_DIR}/postgres/data" "${DATA_DIR}/postgres/cache" \
        "${DATA_DIR}/data" "${DATA_DIR}/faiss" "${DATA_DIR}/locks" \
        "${DATA_DIR}/backups" "${DATA_DIR}/trash" \
        "${LOG_DIR}/postgres" "${DATA_DIR}/versions" \
        "${DATA_DIR}/watch_catalog" "${DATA_DIR}/watched" "$SECRETS"
    seed_watch_catalog_example
    ensure_runtime_config_template
    "${BIN}/casmgr-install-server-config" "$CONFIG_DIR" "$RUNTIME_CONFIG_TEMPLATE" "$CASMGR_GROUP"
    if [[ ! -f "${SECRETS}/.env" ]]; then
        if [[ -f /usr/share/casmgr-server/secrets.env.template ]]; then
            cp /usr/share/casmgr-server/secrets.env.template "${SECRETS}/.env"
        elif [[ -f /usr/share/doc/casmgr-server/secrets.env.template ]]; then
            cp /usr/share/doc/casmgr-server/secrets.env.template "${SECRETS}/.env"
        elif [[ -f /usr/share/doc/casmgr-server/secrets.env.template.gz ]]; then
            zcat /usr/share/doc/casmgr-server/secrets.env.template.gz >"${SECRETS}/.env"
        fi
        if grep -q CHANGE_ME "${SECRETS}/.env" 2>/dev/null; then
            super_pw="$(openssl rand -base64 24 | tr -d '/+=' | head -c 32)"
            app_pw="$(openssl rand -base64 24 | tr -d '/+=' | head -c 32)"
            sed -i "s/CHANGE_ME_SUPER/${super_pw}/" "${SECRETS}/.env"
            sed -i "s/CHANGE_ME_APP/${app_pw}/" "${SECRETS}/.env"
            echo "Generated PostgreSQL passwords in ${SECRETS}/.env"
        fi
        chown "$CASMGR_USER:$CASMGR_GROUP" "${SECRETS}/.env"
        chmod 0640 "${SECRETS}/.env"
    fi
    chown -R 999:999 "${DATA_DIR}/postgres/data" "${DATA_DIR}/postgres/cache" 2>/dev/null || true
}

fix_mtls_permissions() {
    local mtls="${CONFIG_DIR}/mtls"
    [[ -d "$mtls" ]] || return 0
    chown root:"$CASMGR_GROUP" "$mtls"
    chmod 0750 "$mtls"
    local f
    for f in "$mtls"/*; do
        [[ -e "$f" ]] || continue
        chown root:"$CASMGR_GROUP" "$f"
        chmod 0640 "$f"
    done
}

_casmgr_pinned_image_ref() {
    local version="latest"
    if [[ -f "${CONFIG_DIR}/.env" ]]; then
        version="$(grep -m1 '^CASMGR_VERSION=' "${CONFIG_DIR}/.env" | cut -d= -f2-)"
        [[ -n "$version" ]] || version="latest"
    fi
    echo "vasilyvz/casmgr:${version}"
}

provision_image() {
    local tarball="/usr/share/casmgr/casmgr-image.tar.gz"
    if [[ -f "$tarball" ]] && command -v docker >/dev/null 2>&1; then
        echo "Loading casmgr image from ${tarball} ..."
        docker load -i "$tarball" || echo "WARNING: docker load failed for ${tarball}" >&2
    fi
    # If no tarball was shipped, casmgr.service's ExecStartPre pulls the
    # pinned image (see packaging/systemd/casmgr.service); offline hosts
    # without network access must pre-load the image manually.
}

ensure_instance_uuid() {
    # Best-effort revalidation: casmgr-install-server-config already assigns
    # a fresh UUID4 via uuidgen on first install (REPLACE_ON_INSTALL case, no
    # Python required). This extra pass runs the package's stricter UUID4
    # check/repair *inside* the image, since no host venv exists anymore.
    [[ -f "$CONFIG" ]] || return 0
    command -v docker >/dev/null 2>&1 || return 0
    local image
    image="$(_casmgr_pinned_image_ref)"
    docker image inspect "$image" >/dev/null 2>&1 || return 0
    docker run --rm -v "${CONFIG_DIR}:${CONFIG_DIR}" "$image" \
        /opt/casmgr/venv/bin/python -c "
from code_analysis.core.config_instance_uuid import ensure_instance_uuid_in_config
ensure_instance_uuid_in_config('${CONFIG}')
" 2>/dev/null || true
}

config_ready_for_server() {
    [[ -f "$CONFIG" ]] || return 1
    if grep -qE 'CHANGE_ME|MCP_PROXY_HOST|CHUNKER_HOST|EMBEDDING_HOST' "$CONFIG" 2>/dev/null; then
        return 1
    fi
    local cert
    for cert in server.crt server.key ca.crt client.crt client.key; do
        [[ -f "${CONFIG_DIR}/mtls/${cert}" ]] || return 1
    done
    return 0
}

migrate_legacy_service() {
    if systemctl list-unit-files code-analysis-server.service >/dev/null 2>&1; then
        systemctl disable code-analysis-server.service 2>/dev/null || true
        systemctl reset-failed code-analysis-server.service 2>/dev/null || true
    fi
    local legacy_svc
    for legacy_svc in casmgr-server.service casmgr-postgres.service; do
        if systemctl list-unit-files "$legacy_svc" >/dev/null 2>&1; then
            systemctl stop "$legacy_svc" 2>/dev/null || true
            systemctl disable "$legacy_svc" 2>/dev/null || true
            systemctl reset-failed "$legacy_svc" 2>/dev/null || true
        fi
    done
    if [[ -d /usr/lib/code-analysis-server && ! -e /usr/lib/casmgr-server/.migrated-from-code-analysis-server ]]; then
        echo "Note: legacy /usr/lib/code-analysis-server detected; casmgr-server uses /usr/lib/casmgr-server."
    fi
}

print_post_install_banner() {
    local container_started=0
    local pw_status="review and back up"
    local config_pending=()
    local cert missing

    if config_ready_for_server; then
        container_started=1
    fi
    if [[ -f "${SECRETS}/.env" ]]; then
        if grep -qE 'CHANGE_ME' "${SECRETS}/.env" 2>/dev/null; then
            pw_status="REPLACE CHANGE_ME_* placeholders, then run: sudo casmgr-pg-set-password"
        else
            pw_status="auto-generated on first install — back up securely; rotate with: sudo casmgr-pg-set-password"
        fi
    else
        pw_status="create ${SECRETS}/.env from secrets.env.template"
    fi
    if [[ -f "$CONFIG" ]]; then
        grep -q 'CHANGE_ME' "$CONFIG" 2>/dev/null && config_pending+=("server.advertised_host (CHANGE_ME)")
        grep -q 'MCP_PROXY_HOST' "$CONFIG" 2>/dev/null && config_pending+=("registration URLs (MCP_PROXY_HOST)")
        grep -q 'CHUNKER_HOST' "$CONFIG" 2>/dev/null && config_pending+=("code_analysis.chunker.url (CHUNKER_HOST, port 8009)")
        grep -q 'EMBEDDING_HOST' "$CONFIG" 2>/dev/null && config_pending+=("code_analysis.embedding.host (EMBEDDING_HOST, port 8001)")
        grep -q 'REPLACE_ON_INSTALL' "$CONFIG" 2>/dev/null && config_pending+=("registration.instance_uuid (REPLACE_ON_INSTALL)")
    fi
    for cert in server.crt server.key ca.crt client.crt client.key; do
        [[ -f "${CONFIG_DIR}/mtls/${cert}" ]] || missing=1
    done

    echo ""
    echo "================================================================================"
    echo "casmgr-server installed — configure before production use"
    echo "================================================================================"
    echo ""
    echo "Paths:"
    echo "  Config:  ${CONFIG}"
    if [[ -f "$CONFIG" && -f "$CONFIG_TEMPLATE_LOCAL" ]]; then
        echo "  Template (package): ${CONFIG_TEMPLATE_LOCAL}"
    fi
    echo "  Compose: ${COMPOSE_FILE}"
    echo "  Secrets: ${SECRETS}/.env"
    echo "  mTLS:    ${CONFIG_DIR}/mtls/"
    echo "  Logs:    ${LOG_DIR} (host) and 'sudo casmgr logs' / 'docker logs casmgr'"
    echo "  Docs:    man casmgr-server   info casmgr-server"
    echo ""
    echo "1) Edit ${CONFIG}:"
    if ((${#config_pending[@]} > 0)); then
        local item
        for item in "${config_pending[@]}"; do
            echo "     - ${item}"
        done
    else
        echo "     - registration URLs, advertised_host, watch_dirs (if not set yet)"
    fi
    echo "     - code_analysis.chunker.url and embedding.host (inside code_analysis, not root;"
    echo "       SVO chunker :8009, embedding API :8001 — required for vectorization worker)"
    echo "     - MCP port default: 15010 (dev git checkout uses 15000 on the same host)"
    echo ""
    echo "2) PostgreSQL passwords (${SECRETS}/.env):"
    echo "     ${pw_status}"
    echo "     Apply after any edit:"
    echo "       sudo casmgr-pg-set-password"
    echo "     (first-time role/database creation happens automatically on container boot)"
    echo ""
    echo "3) mTLS certificates under ${CONFIG_DIR}/mtls/:"
    if [[ -n "${missing:-}" ]]; then
        echo "     MISSING: server.crt, server.key, ca.crt, client.crt, client.key"
    else
        echo "     OK: required files present"
    fi
    echo ""
    if (( container_started == 1 )); then
        echo "4) casmgr.service was started (config and mTLS look ready)."
        echo "     sudo systemctl status casmgr.service"
    else
        echo "4) casmgr.service was started, but config/mTLS still need attention (see 1-3)."
        echo "     After fixing them: sudo systemctl restart casmgr.service"
    fi
    echo ""
    echo "Password rotation: edit ${SECRETS}/.env, then sudo casmgr-pg-set-password"
    echo "  (set POSTGRES_SUPERUSER_CONNECT_PASSWORD to the old superuser password when rotating)"
    echo ""
    echo "================================================================================"
    echo ""
}

if [[ "$1" == "configure" ]]; then
    create_user
    create_layout
    fix_mtls_permissions
    provision_image
    migrate_legacy_service
    systemctl daemon-reload
    systemctl enable casmgr.service || true
    if command -v docker >/dev/null 2>&1; then
        systemctl start casmgr.service || true
        ensure_instance_uuid
    fi
    print_post_install_banner
fi

exit 0
