# debian/prerm
#!/bin/bash
# debian/prerm — stop services on upgrade/remove.

set -e

stop_unit() {
    local svc="$1"
    if systemctl is-active --quiet "$svc" 2>/dev/null; then
        systemctl stop "$svc" || true
    fi
    if systemctl is-enabled --quiet "$svc" 2>/dev/null; then
        systemctl disable "$svc" || true
    fi
}

case "$1" in
    remove|upgrade|deconfigure)
        stop_unit casmgr
        stop_unit casmgr-server
        stop_unit casmgr-postgres
        if command -v docker >/dev/null 2>&1 && [[ -f /etc/casmgr/docker-compose.yml ]]; then
            docker compose -f /etc/casmgr/docker-compose.yml down 2>/dev/null || true
        fi
        ;;
esac

exit 0
