# debian/prerm
#!/bin/bash
# debian/prerm - pre-removal script

set -e

SERVICE_NAME="code-analysis-server"

# Stop service before removal
if systemctl is-active --quiet "$SERVICE_NAME" 2>/dev/null; then
    echo "Stopping $SERVICE_NAME service..."
    systemctl stop "$SERVICE_NAME" || true
fi

# Disable service
if systemctl is-enabled --quiet "$SERVICE_NAME" 2>/dev/null; then
    systemctl disable "$SERVICE_NAME" || true
fi

exit 0

