# debian/prerm
#!/bin/bash
# debian/prerm — stop services; remove PostgreSQL container on package remove.

set -e

BIN="/usr/lib/casmgr/bin"

stop_unit() {
    local svc="$1"
    if systemctl is-active --quiet "$svc" 2>/dev/null; then
        systemctl stop "$svc" || true
    fi
    if systemctl is-enabled --quiet "$svc" 2>/dev/null; then
        systemctl disable "$svc" || true
    fi
}

case "$1" in
    remove|upgrade|deconfigure)
        stop_unit casmgr-server
        stop_unit casmgr-postgres
        if [[ -x "${BIN}/casmgr-postgres-container" ]]; then
            "${BIN}/casmgr-postgres-container" remove || true
        fi
        ;;
esac

exit 0
