# docker/casmgr/s6/s6-rc.d/40-casmgr/run
#!/command/with-contenv bash
# Longrun (root -> gosu casuser): the code_analysis MCP daemon and all
# in-process workers (indexing / vectorization / file-watcher / queue-manager).
# Drops privileges to casuser (CASMGR_UID:CASMGR_GID, created by 05-init-user)
# so files the daemon creates in the bind-mounted watched trees are
# casuser-owned on the host, not root:root. Must run via the venv interpreter
# so sys.executable-based quality-tool resolution (flake8/mypy/black/isort/
# bandit) finds them in /opt/casmgr/venv.
set -uo pipefail

: "${CASMGR_UID:?CASMGR_UID not set}"
: "${CASMGR_GID:?CASMGR_GID not set}"

if [[ -f /var/casmgr/secrets/.env ]]; then
    set -a
    # shellcheck disable=SC1091
    source /var/casmgr/secrets/.env
    set +a
fi

# Guard against a fresh-first-boot race: 05-init-user's chown of /var/log/casmgr
# runs in parallel with 10-pg-initdb/20-postgres, before server_startup.log is
# ever created; that file is actually first created (root-owned, 644) later, by
# 30-pg-init importing code_analysis (whose package __init__ installs the log
# dir) as root. Without this, casuser can read but not append to it, so the
# daemon's own append_server_startup_log() calls fail silently. Re-chown here,
# immediately before the privilege drop, so this is correct on every boot
# (fresh volume or restart) regardless of what any earlier root step created.
chown -R "${CASMGR_UID}:${CASMGR_GID}" /var/log/casmgr

# Stable CWD: s6 launches longruns with CWD = the s6-rc service dir
# (/run/s6-rc:.../servicedirs/40-casmgr). Queue-manager worker subprocesses
# inherit that CWD, and a code path that resolves "config.json" relatively then
# fails ("No such file or directory: .../40-casmgr/config.json"), breaking every
# use_queue command (git_pull/fetch/push/clone, long analyses). Anchor CWD at
# /etc/casmgr so a relative config.json resolves to the real config; the config's
# own storage paths are absolute under /var/casmgr, so nothing else is affected.
cd /etc/casmgr || exit 1

exec gosu "${CASMGR_UID}:${CASMGR_GID}" /opt/casmgr/venv/bin/python -m code_analysis.main --config /etc/casmgr/config.json --foreground
